# docker.conf
description "Docker daemon"

start on filesystem and started lxc-net
stop on runlevel [!2345]

respawn

script
        DOCKER=/usr/bin/$UPSTART_JOB
        DOCKER_OPTS=
        if [ -f /etc/default/$UPSTART_JOB ]; then
                . /etc/default/$UPSTART_JOB
        fi
        "$DOCKER" $DOCKER_OPTS -b bridge0 -dns 172.30.1.1 -H tcp://0.0.0.0:4243 -H unix:///var/run/docker.sock -d
end script
